# curator4-1896754666.conf
[program:curator4-1896754666]
command=python3 /home/curator4/curator4-1896754666/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1896754666/test.err.log
stdout_logfile=/home/curator4/curator4-1896754666/test.out.log